# deploy/systemd/mim-objective75-nightly-summary.timer
[Unit]
Description=Run Objective 75 Nightly Summary at 00:20

[Timer]
OnCalendar=*-*-* 00:20:00
Persistent=true
RandomizedDelaySec=120
Unit=mim-objective75-nightly-summary.service

[Install]
WantedBy=timers.target
