# deploy/systemd/mim-watch-objective75-cycle-pass.service
[Unit]
Description=MIM Watch Objective 75 Cycle PASS Staleness
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=testpilot
Group=testpilot
WorkingDirectory=/home/testpilot/mim
Environment=MAX_IDLE_SECONDS=1200
Environment=POLL_SECONDS=30
Environment=COOLDOWN_SECONDS=300
ExecStart=/home/testpilot/mim/scripts/watch_objective75_cycle_pass.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
