# deploy/systemd/mim-watch-tod-liveness.service
[Unit]
Description=MIM Watch TOD Liveness
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=testpilot
Group=testpilot
WorkingDirectory=/home/testpilot/mim
Environment=POLL_SECONDS=3
Environment=STALE_SECONDS=45
Environment=COOLDOWN_SECONDS=30
Environment=TRIGGER_PROTECT_SECONDS=20
ExecStart=/home/testpilot/mim/scripts/watch_tod_liveness.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
