# deploy/systemd-user/mim-handoff-watcher.service
[Unit]
Description=MIM Handoff Watcher
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/testpilot/mim
Environment=MIM_HANDOFF_ROOT=/home/testpilot/mim/handoff
Environment=MIM_SHARED_ROOT=/home/testpilot/mim/runtime/shared
Environment=MIM_HANDOFF_POLL_INTERVAL_SECONDS=2
ExecStart=/home/testpilot/mim/.venv/bin/python /home/testpilot/mim/scripts/watch_handoff_inbox.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target