# deploy/systemd-user/mim-watch-mim-coordination-responder.service
[Unit]
Description=MIM Watch Coordination Responder (User)
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/testpilot/mim
Environment=POLL_SECONDS=3
ExecStart=/home/testpilot/mim/scripts/watch_mim_coordination_responder.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
