# deploy/systemd-user/mim-desktop-shell.service
[Unit]
Description=MIM Desktop Shell
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/testpilot/mim
Environment=MIM_UI_URL=http://127.0.0.1:8001/mim
Environment=MIM_ELECTRON_DISABLE_SANDBOX=1
ExecCondition=/usr/bin/env bash -lc 'command -v npm >/dev/null 2>&1'
ExecStart=/usr/bin/env bash /home/testpilot/mim/scripts/run_mim_desktop_shell.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
