# deploy/systemd-user/mim-objective75-overnight-loop.service
[Unit]
Description=MIM Objective 75 Overnight Loop (User)
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/testpilot/mim
Environment=LOOP_INTERVAL_SECONDS=120
Environment=SOAK_WAIT_TIMEOUT_SECONDS=120
Environment=SOAK_CYCLES=2
Environment=SOAK_SLEEP_SECONDS=2
Environment=REQUIRE_TOD_RESULT_COMPLETED=1
ExecStart=/home/testpilot/mim/scripts/run_objective75_overnight_loop.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
